# server_voice_assistant/voice-doubao.service
[Unit]
Description=Doubao Realtime Voice Assistant
After=network.target sound.target

[Service]
Type=simple
User=ladydd
Group=audio
WorkingDirectory=/home/ladydd/voice-doubao
ExecStart=/home/ladydd/.hermes/hermes-agent/venv/bin/python /home/ladydd/voice-doubao/main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
